# rootfs/lib/udev/rules.d/51-automedia.rules
ACTION=="change", SUBSYSTEM=="block", RUN+="/opt/arm/scripts/arm_wrapper.sh %k"