# rootfs/lib/udev/rules.d/51-automedia.rules
ACTION=="change", SUBSYSTEM=="block", RUN+="python3 /opt/arm/arm/ripper/main.py -d %k"