# rootfs/lib/udev/rules.d/51-automedia.rules
ACTION=="change", SUBSYSTEM=="block", RUN+="/etc/arm_wrapper.sh %k"