# rootfs/etc/services.d/logserver/run
#!/usr/bin/with-contenv sh

cd /opt/arm
exec python3 runui.py
