# rootfs/etc/services.d/logserver/run
#!/usr/bin/with-contenv sh

cd /opt/arm
exec s6-setuidgid arm python3 runui.py
