# rootfs/etc/services.d/udev-log/run
#!/usr/bin/with-contenv sh

exec /lib/systemd/systemd-udevd --debug